# conf/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ── Cloudflare WARP daemon ────────────────────────────────────────────────────
[program:warp-svc]
command=/usr/bin/warp-svc
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/warp-svc.log
stderr_logfile=/var/log/warp-svc-err.log
startsecs=5

# ── WARP initial setup (runs once after warp-svc starts) ─────────────────────
[program:warp-setup]
command=/warp-setup.sh
autostart=true
autorestart=false
priority=20
startsecs=0
startretries=3
stdout_logfile=/var/log/warp-setup.log
stderr_logfile=/var/log/warp-setup-err.log

# ── Dante SOCKS5 proxy ────────────────────────────────────────────────────────
[program:danted]
command=/usr/sbin/danted -f /etc/danted.conf
autostart=false
autorestart=true
priority=30
stdout_logfile=/var/log/danted.log
stderr_logfile=/var/log/danted-err.log
